# totrans/sec-mon-wazuh_06.md.json
["```\nJan 21 15:15:45 myCP Checkpoint: 21Jan2019 15:15:45 monitor 10.0.10.1 <bond0 Protection Name:Header Rejection;Severity:4;Confidence Level:4;protection_id:HttpHeaderRejection;SmartDefense Profile:SU2_Protection;Performance Impact:2;Industry Reference:CVE-2002-0032, CAN-2003-0237, CAN-2002-0254, CVE-2002-0155, CAN-2003-0397, CAN-2002-0314;Protection Type:protection;Signature Info:^User-Agent[^I ]*:[^I ]*.*esb|ESB;Update Version:634182243;rule:26;rule_uid:{405CB782-3274-4D7F-8AAA-4FB24CE726A0};resource:<http://dnl-02.geo.kaspersky.com/bases/av/kdb/i386/kdb-i386-1211g.xml.klz;reject_id:5accf7c4-10053-c00080a-c0000003;web_client_type:Other:> *BcfBAAAAgCCAAEFBAAwQfKXVzrzGvyfPESboPxow0mHhxRLAXAQAAIAAKAA=;Attack Info:WSE0100001 header rejection pattern found in request;attack:Header Rejection;src:10.20.10.1;dst:1.1.1.1;proto:6;proxy_src_ip:10.10.10.1;product:SmartDefense;service:80;s_port:51642;FollowUp:Not Followed;product_family:Network\n```", "```\n<decoder name=\"checkpoint-syslog\">\n  <program_name>^Checkpoint</program_name>\n  <prematch>^\\\\s*\\\\S+ \\\\d\\\\d:\\\\d\\\\d:\\\\d\\\\d </prematch>\n</decoder>\n```", "```\n<decoder name=\"checkpoint-syslog-fw\">\n  <parent>checkpoint-syslog</parent>\n  <type>firewall</type>\n  <prematch offset=\"after_parent\">^drop|^accept|^reject</prematch>\n  <regex offset=\"after_parent\">^(\\\\w+)\\\\s+\\\\S+ \\\\p\\\\S+ rule:\\\\.+</regex>\n  <regex>src: (\\\\S+); dst: (\\\\S+); proto: (\\\\S+);</regex>\n  <order>action,srcip,dstip,protocol</order>\n</decoder>\n<decoder name=\"checkpoint-syslog-fw\">\n  <parent>checkpoint-syslog</parent>\n  <type>firewall</type>\n  <regex offset=\"after_regex\">service: (\\\\d+); s_port: (\\\\d+);</regex>\n  <order>dstport,srcport</order>\n</decoder>\n<decoder name=\"checkpoint-syslog-ids\">\n  <parent>checkpoint-syslog</parent>\n  <type>ids</type>\n  <prematch offset=\"after_parent\">^monitor|^drop</prematch>\n  <regex offset=\"after_prematch\">attack:\\\\s*(\\\\.+);\\\\s*</regex>\n  <regex>src:\\\\s*(\\\\S+);\\\\s*dst:\\\\s*(\\\\S+);\\\\s*</regex>\n  <regex>proto:\\\\s*(\\\\S+);</regex>\n  <order>extra_data, srcip, dstip, protocol</order>\n  <fts>name, extra_data, srcip, dstip</fts>\n  <ftscomment>First time Checkpoint rule fired.</ftscomment>\n</decoder>\n```", "```\n<group name=\"checkpoint-syslog,\">\n  <!--Generic rule -->\n  <rule \"d=\"64\"00\" lev\"l\"\"3\">\n    <decoded_as>checkpoint-syslog</decoded_as>\n    <description>Checkpoint $(type) event</description>\n  </rule>\n```", "```\n    <ossec_config>\n      <global>\n        <jsonout_output>yes</jsonout_output>\n        <alerts_log>yes</alerts_log>\n        <logall>yes</logall>\n        <logall_json>yes</logall_json>\n    </ossec_config>\n    ```", "```\n    filebeat service by editing /etc/filebeat/filebeat.yml and changing archives: value to true.Next, restart the `filebeat` service as follows:\n\n    ```", "```\n\n    ```", "```\n<rule id=\"100009\" level=\"10\" frequency=\"8\" timeframe=\"120\" ignore=\"60\">\n    <if_matched_sid>100001</if_matched_sid>\n    <description>sshd: brute force attack</description>\n    <same_srcip />\n    <mitre>\n      <id>T1110</id>\n    </mitre>\n  </rule>\n```", "```\n    select pid, name, uid, resident_size from processes order by resident_size desc limit 10;\n    ```", "```\n    select count(pid) as total, name from processes group by name order by total desc limit 10;\n    ```", "```\n    apt-key command. This key is necessary to validate the Osquery packages you will be installing:\n\n    ```", "```\n\n    ```", "```\n    add-apt-repository 'deb [arch=amd64] https://pkg.osquery.io/deb> deb main'\n    apt-get install command to install Osquery. This will get Osquery from the newly added repository and install it:\n\n    ```", "```\n\n    ```", "```\n    <!-- Osquery integration -->\n    <wodle name=\"osquery\">\n    <disabled>no</disabled>\n    <run_daemon>yes</run_daemon> <log_path>/var/log/osquery/osqueryd.results.log</log_path> <config_path>/etc/osquery/osquery.conf</config_path> <add_labels>yes</add_labels>\n     </wodle>\n    ```", "```\n    nano /etc/osquery/osquery.conf to view the default packs:\n\n    ```", "```\n\n    *   In the preceding code, we can see the following:*   `osquery-monitoring.conf` is an Osquery configuration file to collect information about every Osquery pack, including general performance and versions*   `incident-response.conf` retrieves information about crontab, the loginwindow process, a list of open sockets, a list of mounted drives, and so on*   `it-compliance.conf` collects information about active directory, the operating system, shared services, browser plugins, Windows drivers, a list of USB drives, and so on*   `vuln-management.conf` retrieves information about installed applications, browser plugins, and Chrome extensions*   `hardware-monitoring.conf` gathers hardware-related information such as PCI devices, fan speed, an inventory of USB drives, kernel modules, and so on*   `ossec-rootkit.conf` collects information about rootkits\n    ```", "```\n    systemctl restart osqueryd\n    ```", "```\nselect command, path from crontab;\n```", "```\nselect name from kernel_modules;\n```", "```\nselect pr.pid, pr.name, usr.username, pr.path, pr.cmdline from processes pr LEFT JOIN users usr ON pr.uid = usr.uid WHERE pr.cmdline != '';\n```", "```\n{\n  \"queries\": {\n    \"Services\": {\n        \"query\": \"SELECT * FROM services WHERE start_type='DEMAND_START' OR start_type='AUTO_START';\",\n        \"interval\": 3600,\n        \"description\": \"Lists all installed services configured to start automatically at boot - ATT&CK T1050\",\n        \"removed\": false\n    },\n    \"Snapshot_services\": {\n        \"query\": \"SELECT * FROM services;\",\n        \"interval\": 28800,\n        \"description\": \"Snapshot Services query\",\n        \"snapshot\": true\n    },\n      \"OtherServices\": {\n        \"query\": \"SELECT name, display_name, status, start_type, path, module_path FROM services;\",\n        \"interval\": 3600,\n        \"description\": \"Services whose executables are placed in unfamiliar folders- ATT&CK T1543.003\",\n        \"removed\": false\n    }\n   }\n}\n```", "```\n <wodle name=\"command\">\n    <disabled>no</disabled>\n    <tag>tasklist</tag>\n    <command>PowerShell.exe C:\\\\activeTasks.bat</command>\n    <interval>2m</interval>\n    <run_on_start>yes</run_on_start>\n    <timeout>10</timeout>\n  </wodle>\n```", "```\n<localfile>\n<log_format>full_command</log_format> <command><COMMAND></command>\n<frequency>120</frequency>\n </localfile>\n```", "```\n    ifconfig, netstat, route, arp, rarp, and so on.\n    ```", "```\n    <ossec_config>\n    <localfile>\n     <log_format>full_command</log_format>\n    <command>netstat -tulpn</command>\n    <alias>netstat listening ports</alias>\n     <frequency>360</frequency>\n     </localfile>\n    <log_format>full_command</log_format>: This specifies the log format. In this case, it is set to full_command, indicating that the log consists of the full output.\n    ```", "```\n    systemctl restart wazuh-agent\n    ```", "```\n<rule id=\"533\" level=\"7\">\n<if_sid>530</if_sid>\n<match>ossec: output: 'netstat listening ports</match>\n<check_diff />\n<description>Listened ports status (netstat) changed (new port opened or closed).</description>\n<group>pci_dss_10.2.7,pci_dss_10.6.1,gpg13_10.1,gdpr_IV_35.7.d,hipaa_164.312.b,nist_800_53_AU.14,nist_800_53_AU.6,tsc_CC6.8,tsc_CC7.2,tsc_CC7.3,</group>\n</rule>\n```", "```\n<group name=\"ossec,\">\n<rule id=\"500\" level=\"0\">\n<category>ossec</category>\n<decoded_as>ossec</decoded_as>\n<description>Grouping of ossec rules.</description>\n</rule>\n</group>\n```"]